# supervisord-docker.ini
[program:docker]
command=/usr/local/bin/dockerd-entrypoint.sh
autorestart=true
